# ops/bluette.service.example
[Unit]
Description=Bluette Server
After=network.target

[Service]
Type=simple
User=bluette
Group=bluette
WorkingDirectory=/opt/bluette_server
EnvironmentFile=/etc/bluette/bluette.env
ExecStart=/opt/bluette_server/_build/prod/rel/bluette_server/bin/bluette_server start
Restart=always
RestartSec=5
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target